# deployment/supervisor.conf
; Supervisor configuration for technical_ai on VPS
; Place in /etc/supervisor/conf.d/technical_ai.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update

;----------------------------------------------------------
; TRAINING QUEUE - Heavy document processing (CPU/memory)
; Runs on dedicated workers, isolated from user-facing chat
;----------------------------------------------------------
[program:technical_ai-training]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/technical_ai/artisan queue:work database --queue=training --sleep=3 --tries=3 --timeout=300 --memory=1024
autostart=true
autorestart=true
stopastype=TERM
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/technical_ai/storage/logs/training-worker.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=330

;----------------------------------------------------------
; DEFAULT QUEUE - Lightweight tasks (emails, notifications)
;----------------------------------------------------------
[program:technical_ai-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/technical_ai/artisan queue:work database --queue=default --sleep=3 --tries=3 --timeout=60 --memory=256
autostart=true
autorestart=true
stopastype=TERM
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/technical_ai/storage/logs/default-worker.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=70

;----------------------------------------------------------
; GROUP - Manage all workers together
;----------------------------------------------------------
[group:technical_ai-workers]
programs=technical_ai-training,technical_ai-default
priority=999
